# remote-trigger.service
[Unit]
Description=remote-trigger

[Service]
ExecStart=/usr/local/bin/remote-trigger
Restart=on-abort

[Install]
WantedBy=multi-user.target